# systemd/reclaim-rebates.service
[Unit]
Description=Reclaim Rebate Distribution
After=network-online.target reclaim-rpc.service

[Service]
Type=oneshot
WorkingDirectory=/root/reclaim-fi
EnvironmentFile=/opt/mev/.env
ExecStart=/root/reclaim-fi/venv/bin/python rebate_service.py
StandardOutput=journal
StandardError=journal
