# systemd/reclaim-rebates.timer
[Unit]
Description=Reclaim Rebate Distribution Timer (hourly)

[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
